# deploy/systemd/shennong-llama-cpp-proxy.service
[Unit]
Description=Proxy Shennong OS Agent requests to loopback-only llama.cpp adapter
Requires=shennong-llama-cpp-proxy.socket
After=pi-llama-reasoning-proxy.service

[Service]
ExecStart=/lib/systemd/systemd-socket-proxyd --connections-max=32 127.0.0.1:8081
DynamicUser=true
PrivateDevices=true
PrivateTmp=true
ProtectControlGroups=true
ProtectHome=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
NoNewPrivileges=true
RestrictAddressFamilies=AF_INET AF_INET6
